# deploy/config/systemd/vault-web.service
[Unit]
Description=ARCHITECT // VAULT - Web Application
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=notify
User=vault
Group=vault
WorkingDirectory=/opt/vault

# Load environment variables
EnvironmentFile=/etc/vault/environment

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/vault /var/log/vault
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

# Start command
ExecStart=/opt/vault/venv/bin/uvicorn app.main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 4 \
    --log-config /etc/vault/logging.conf \
    --proxy-headers \
    --forwarded-allow-ips '*'

# Restart policy
Restart=always
RestartSec=10
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vault-web

[Install]
WantedBy=multi-user.target
